# fpm-nginx/s6-rc.d/nginx/run
#!/command/with-contenv sh
# shellcheck shell=sh
exec 2>&1
# nginx cannot expand env vars in its config, so render the template into /run
# (writable by www-data) and point nginx at it. Only our ${SERVER_ROOT} and
# ${FASTCGI_READ_TIMEOUT} placeholders are substituted; nginx's own $uri/$document_root
# are left intact. SERVER_TIMEOUT (an image ENV, default 30) is nginx's fastcgi read
# timeout, so nginx doesn't 504 before PHP finishes - raise it for long requests.
# Must be a positive integer of seconds: nginx reads 0 as "time out immediately", so a
# 0 (or non-numeric) value is a misconfiguration, not "unlimited".
case "$SERVER_TIMEOUT" in
    '' | *[!0-9]* | 0)
        echo "nginx: SERVER_TIMEOUT must be a positive integer of seconds (got '${SERVER_TIMEOUT:-}')." >&2
        exit 1 ;;
esac
sed -e "s#\${SERVER_ROOT}#${SERVER_ROOT}#g" \
    -e "s#\${FASTCGI_READ_TIMEOUT}#${SERVER_TIMEOUT}#g" \
    /etc/nginx/nginx.conf.template > /run/nginx.conf
exec nginx -c /run/nginx.conf -g "daemon off;"
