# fpm-nginx/s6-rc.d/nginx/run
#!/command/with-contenv sh
# shellcheck shell=sh
exec 2>&1
exec nginx -g "daemon off;"
