# fpm-nginx/s6-rc.d/nginx/run
#!/command/with-contenv sh
# shellcheck shell=sh
exec 2>&1
# nginx cannot expand env vars in its config, so render ${SERVER_ROOT} from the
# template into /run (writable by www-data) and point nginx at it. Only
# ${SERVER_ROOT} is substituted; nginx's own $uri/$document_root are left intact.
sed "s#\${SERVER_ROOT}#${SERVER_ROOT}#g" /etc/nginx/nginx.conf.template > /run/nginx.conf
exec nginx -c /run/nginx.conf -g "daemon off;"
