# fpm-apache/s6-rc.d/apache/run
#!/command/with-contenv sh
# shellcheck shell=sh
exec 2>&1
# vhost.conf's ProxyTimeout reads ${SERVER_TIMEOUT} (an image ENV, default 30, already
# in httpd's environment); it makes Apache wait as long as PHP may run so it doesn't 504
# before PHP finishes - raise it for long requests. Must be a positive integer of seconds:
# Apache's ProxyTimeout must be >= 1, so a 0 (or non-numeric) value is a misconfiguration.
case "$SERVER_TIMEOUT" in
    '' | *[!0-9]* | 0)
        echo "apache: SERVER_TIMEOUT must be a positive integer of seconds (got '${SERVER_TIMEOUT:-}')." >&2
        exit 1 ;;
esac
# Run Apache in the foreground regardless of distro: Debian ships apache2ctl
# (which sources /etc/apache2/envvars), Alpine ships httpd.
if command -v apache2ctl >/dev/null 2>&1; then
    exec apache2ctl -D FOREGROUND
fi
exec httpd -D FOREGROUND
