# fpm-apache/s6-rc.d/apache/run
#!/command/with-contenv sh
exec 2>&1
# Run Apache in the foreground regardless of distro: Debian ships apache2ctl
# (which sources /etc/apache2/envvars), Alpine ships httpd.
if command -v apache2ctl >/dev/null 2>&1; then
    exec apache2ctl -D FOREGROUND
fi
exec httpd -D FOREGROUND
